# systemd/plymouth-killer.service
[Unit]
Description=Kill plymouthd left over from initramfs
# ophub 预编译的 uInitrd 内置了 plymouthd，切换到真实 rootfs 后进程残留
# 本机没有完整 plymouth 运行时，plymouthd 会持续占用 CPU
# 此 service 在启动完成后将其 kill 掉
After=multi-user.target
ConditionPathExists=/run/plymouth/pid

[Service]
Type=oneshot
ExecStart=/bin/bash -c '\
  PID=$(cat /run/plymouth/pid 2>/dev/null); \
  [ -n "$PID" ] && kill "$PID" 2>/dev/null; \
  rm -f /run/plymouth/pid; \
  /usr/lib/debian/fix-armbianenv 2>/dev/null || true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
